# docker/supervisord.conf
; BetterDesk — Supervisord Configuration
; Manages both Go server and Node.js console in a single container

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ---- BetterDesk Go Server ----
; Signal (21116 TCP/UDP), Relay (21117), API (21114), NAT (21115), WS (21118/21119)
; DB_URL env var is set by entrypoint.sh: SQLite path or postgres:// URI
[program:betterdesk-server]
command=/usr/local/bin/betterdesk-server -mode all -key-file /opt/rustdesk/id_ed25519
directory=/opt/rustdesk
user=betterdesk
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
; Log to Docker stdout/stderr so errors appear in `docker logs`
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SIGNAL_PORT="21116",ENCRYPTED_ONLY="%(ENV_ENCRYPTED_ONLY)s",DB_URL="%(ENV_DB_URL)s",RELAY_SERVERS="%(ENV_RELAY_SERVERS)s"
priority=100

; ---- BetterDesk Node.js Console ----
; Web console (5000), RustDesk Client API (21121)
[program:betterdesk-console]
command=node /app/server.js
directory=/app
user=betterdesk
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=10
; Log to Docker stdout/stderr so errors appear in `docker logs`
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Ensure Client API binds to all interfaces inside container
environment=API_HOST="0.0.0.0",HOST="0.0.0.0"
priority=200
